# build/demo-launcher/supervisord.conf
[supervisord]
nodaemon=true
user=root
#loglevel=error

[inet_http_server]
port = 127.0.0.1:9001
username = dummy
password = dummy


[program:launcher]
command = python3 -Xfrozen_modules=off -m debugpy --listen 0.0.0.0:5678 -m uvicorn launcher.main:app --proxy-headers --host 0.0.0.0 --port 8080 --log-level info --reload-dir /app/launcher
redirect_stderr = true
redirect_stdout = true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=1


[program:firecrest]
command = uvicorn firecrest.main:app --proxy-headers --host 0.0.0.0 --port 5000
redirect_stderr = true
redirect_stdout = true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
autostart=false

[program:firecrest-ui]
directory=/app/web-ui/
command = /bin/bash -c "KEYCLOAK_CLIENT_ID=$(cat /app/config/webui-client) yarn start"   
redirect_stderr = true
redirect_stdout = true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
autostart=false

[supervisorctl]
serverurl=http://127.0.0.1:9001
username = dummy
password = dummy

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface


