# build/docker/pbs-cluster/supervisord.conf
[supervisord]
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -D
priority=10
autostart=true
autorestart=true

[program:init-script]
command=/usr/local/bin/init_sequence.sh
priority=20
autostart=true
autorestart=false
