# build/docker/pbs-cluster/supervisord.conf
[supervisord]
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -D
priority=10
autostart=true
autorestart=true

[program:postgresql]
command=/usr/lib/postgresql/14/bin/postgres -D /var/lib/postgresql/14/main -c config_file=/etc/postgresql/14/main/postgresql.conf
priority=10
user=postgres
autostart=true
autorestart=true
stderr_logfile=/var/log/supervisor/postgresql.err.log
stdout_logfile=/var/log/supervisor/postgresql.out.log

[program:munge]
command=/usr/sbin/munged --foreground
priority=15
user=munge
autostart=true
autorestart=true
stderr_logfile=/var/log/supervisor/munge.err.log
stdout_logfile=/var/log/supervisor/munge.out.log

[program:boot-pbs]
command=/usr/local/bin/boot_pbs.sh
priority=20
stderr_logfile=/var/log/supervisor/init.err.log
stdout_logfile=/var/log/supervisor/init.out.log