# build/demo-sidecar/supervisord.conf
[supervisord]
nodaemon=true
user=root

[inet_http_server]
port = 127.0.0.1:9001


[program:sidecar]
command = uvicorn sidecar.main:app --proxy-headers --host 0.0.0.0 --port 8080
redirect_stderr = true
redirect_stdout = true
stdout_logfile = /app/sidecar.log
priority=1

[program:firecrest]
command = uvicorn firecrest.main:app --proxy-headers --host 0.0.0.0 --port 5000
redirect_stderr = true
redirect_stdout = true
stdout_logfile = /app/firecrest.log
autostart=false



[supervisorctl]
serverurl=http://127.0.0.1:9001

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface


