# build/docker/slurm-cluster/supervisord.conf
[supervisord]
nodaemon=true
user=root

[program:update-jwt-certs]
command = /download_jwks_certificate.sh
stdout_logfile=/var/log/update-jwt-certs.log
user=root
startsecs = 0
autorestart = false
startretries = 1
priority=1

[program:mariadb]
command=/start_db.sh
priority=1

[program:slurmdbd]
command=/start_slurmdbd.sh
stdout_logfile = /var/log/slurm/slurmdbd.log
priority=3

[program:munge]
command=/usr/sbin/munged -F
user=munge
priority=5

[program:sshd]
command=/usr/sbin/sshd -D -e
stdout_logfile=/var/log/sshd_out.log
redirect_stderr=true
priority=6

[program:slurmctld]
command=/start_slurmctld.sh
user=slurm
stdout_logfile = /var/log/slurm/slurmctld.log
priority=10
startsecs=5

[program:slurmd]
command=/usr/sbin/slurmd -D
user=root
stdout_logfile = /var/log/slurm/slurmd.log
priority=16

[program:slurmrestd]
command=/start_slurmrestd.sh
user=slurmapi
priority=16
stdout_logfile = /var/log/slurm/slurmrestd.log
redirect_stderr=true

[supervisorctl]
